# packaging/debian/email-endpoint.service
[Unit]
Description=email-endpoint self-hosted website email endpoint
Documentation=https://github.com/mlahr/email-endpoint
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/default/email-endpoint
ExecStart=/bin/sh -c 'if [ "$EMAIL_ENDPOINT_ENABLED" != "true" ]; then echo "email-endpoint is disabled in /etc/default/email-endpoint"; exit 0; fi; exec /usr/bin/email-endpoint --config "$EMAIL_ENDPOINT_CONFIG"'
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10s
User=email-endpoint
Group=email-endpoint
StateDirectory=email-endpoint
RuntimeDirectory=email-endpoint
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/email-endpoint /run/email-endpoint
SyslogIdentifier=email-endpoint
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
